# anacrontab
# /etc/anacrontab: configuration file for anacron

# See anacron(8) and anacrontab(5) for details.

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin:/home/gholk/.local/bin
HOME=/home/gholk
LOGNAME=gholk

# These replace cron's entries
# 1	5	cron.daily	run-parts --report /etc/cron.daily
# 7	10	cron.weekly	run-parts --report /etc/cron.weekly
# @monthly	15	cron.monthly	run-parts --report /etc/cron.monthly
@monthly      15      backup-firefox backup-firefox-history.sh
1              3      update-feed    /home/gholk/code/feed/script/update-feed.sh

